# src/puppet/modules/contrail-config/templates/contrail-api.ini.erb
[program:contrail-api]
command=/usr/bin/contrail-api --conf_file /etc/contrail/contrail-api.conf --listen_port <%= contrail_api_port_base %>%(process_num)01d --worker_id %(process_num)s
numprocs=<%= contrail_api_nworkers %>
process_name=%(process_num)s
redirect_stderr=true
stdout_logfile= /var/log/contrail/contrail-api-%(process_num)s.log
stderr_logfile=/dev/null
priority=440
autostart=true
killasgroup=true
stopsignal=KILL
