# src/puppet/modules/contrail-config/templates/contrail-discovery.ini.erb
[program:contrail-discovery]
command=/usr/bin/discovery-server --conf_file /etc/contrail/discovery.conf --listen_port <%= contrail_disc_port_base %>%(process_num)01d --worker_id %(process_num)s
numprocs=<%= contrail_disc_nworkers %>
process_name=%(process_num)s
redirect_stderr=true
stdout_logfile= /var/log/contrail/contrail-discovery-%(process_num)s.log
stderr_logfile=/dev/null
priority=430
autostart=true
killasgroup=true
stopsignal=KILL
