# src/puppet/modules/contrail-config/templates/contrail-api.svc.erb
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization API

<%= contrail_supervisorctl_lines %>

#supervisorctl -s unix:///tmp/supervisord_config.sock ${1} `basename ${0}`

